# gdm.upstart
description "Starts the Gnome Display Manager"

start on stopped rc RUNLEVEL=5 and started messagebus

# if you have haldaemon present as well:
#start on pld.sysinit-done and started messagebus and started haldaemon
#stop on pld.shutdown-started
stop on starting rc RUNLEVEL=[!5]

respawn
respawn limit 10 120

exec /usr/sbin/gdm
