# gdm.service
[Unit]
Description=Gnome Display Manager
Requires=dev-tty9.device
After=dev-tty9.device systemd-user-sessions.service

[Service]
ExecStart=/usr/sbin/gdm -nodaemon

[Install]
WantedBy=graphical.target
